# systemd/safe-server-lockdown.timer
[Unit]
Description=Close public SSH @GRACE_MINUTES@ minutes after timer activation (safe-server)

[Timer]
OnActiveSec=@GRACE_MINUTES@min
Unit=safe-server-lockdown.service

[Install]
WantedBy=timers.target
